# supervisor.conf
[program:waldo]
command = /home/deploy/.virtualenvs/waldo/bin/python /sites/waldo/server.py --port=12345
directory = /sites/waldo
user = deploy
autostart = true
autorestart = true

